# crontab/dockerus
@reboot curl http://pyblisher:8000/update/
19 00 * * * curl http://pyblisher:8000/update/
